# home/dot_config/systemd/user/aws-rotate-iam-keys.service.tmpl
[Unit]
Description=Automatically rotate AWS IAM keys

[Service]
Type=oneshot
ExecStart=%h/.scripts/aws-rotate-iam-keys/src/bin/aws-rotate-iam-keys --profile {{ .aws.rotate_keys.profile }}
